# runit/user/sv/wireplumber/run
#!/bin/sh

sv check "$HOME/.runit/runsvdir/pipewire" >/dev/null || exit 1

exec wireplumber
