# runit/user/sv/fcitx5/run
#!/bin/sh

[ -r ./conf ] && . ./conf

exec fcitx5 ${OPTS:-} 2>&1
